# sync/sync-foobar.path
[Unit]
Description=SyncFoobar
[Path]
PathModified=/dev/disk/
PathExists=/dev/disk/by-label/corsair
[Install]
WantedBy=multi-user.target