# sync/sync-foobar.path
[Unit]
Description=SyncFoobar
[Path]
PathModified=/dev/disk/
PathExists=/media/greg/corsair
