# sync/sync-foobar.path
[Unit]
Description=SyncFoobar
[Path]
PathExists=/dev/disk/by-label/corsair
[Install]
WantedBy=multi-user.target