# sync/sync-foobar.service
[Unit]
Description=SyncFoobar
[Service]
ExecStartPre=/bin/sleep 10s
ExecStart=/home/greg/.local/bin/sync-foobar.sh